# oled-display/rootfs/etc/services.d/oled/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting OLED display service..."

# Optionally enable I2C before starting if requested and not already present
if ls /dev/i2c-* >/dev/null 2>&1; then
  bashio::log.info "I2C devices already present ($(ls /dev/i2c-* | xargs)). Skipping enablement."
else
  if bashio::config.true 'system.enable_i2c'; then
    bashio::log.info "system.enable_i2c=true → attempting to enable I2C"
    /usr/bin/enable_i2c.sh || bashio::log.warning "I2C enablement script encountered issues; continuing"
  else
    bashio::log.info "system.enable_i2c not set; skipping I2C enablement"
  fi
fi

python3 /usr/bin/oled_monitor.py