# oled-display/rootfs/etc/services.d/oled/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting OLED display service..."

# Optional HAOS configurators (adamoutler HassOSConfigurator pattern)
if bashio::config.true 'system.enable_ssh'; then
  bashio::log.info "system.enable_ssh=true -> running debug SSH configurator"
  /usr/bin/enable_ssh.sh || bashio::log.warning "SSH configurator encountered issues; continuing"
fi

if bashio::config.true 'system.enable_i2c'; then
  bashio::log.info "system.enable_i2c=true -> running I2C configurator"
  /usr/bin/enable_i2c.sh || bashio::log.warning "I2C configurator encountered issues; continuing"
elif ls /dev/i2c-* >/dev/null 2>&1; then
  bashio::log.info "I2C devices already present ($(ls /dev/i2c-* | xargs))"
else
  bashio::log.info "system.enable_i2c=false and no /dev/i2c-* nodes; skipping I2C configurator"
fi

python3 /usr/bin/oled_monitor.py
